# packaging/vk-turn-proxy-client.service
[Unit]
Description=VK TURN Proxy Client
After=network-online.target wireguard.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/vk-turn-proxy-client
ExecStart=/bin/sh -c '/usr/bin/vk-turn-proxy-client $OPTS | /bin/bash /usr/share/vk-turn-proxy/routes.sh'
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
